# autodev/superagent-auto-improve@orchestrator-dev.timer
# Install at /etc/systemd/system/superagent-auto-improve@orchestrator-dev.timer
# Model-orchestration layer: less frequent, more architecturally sensitive.
[Unit]
Description=Run /orchestrator-dev every 45 minutes

[Timer]
OnBootSec=15min
OnUnitActiveSec=45min
Persistent=true

[Install]
WantedBy=timers.target
